# pod_forgejo/templates/forgejo_update-authorized_keys.path
{{ansible_managed|comment(decoration='#')}}

[Unit]
Description=Update forgejo authorized_keys on file change

[Path]
PathModified={{file_dst_forgejo_ssh_dir}}/authorized_keys

[Install]
WantedBy=multi-user.target
